# kibana.init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          kibana
# Required-Start:    $local_fs $network $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: A browser based interface for Logstash and ElasticSearch
# Description:       This file should be used to start and stop kibana.
### END INIT INFO

PATH=/usr/local/bin:/usr/local/sbin:/bin:/sbin:/usr/bin:/usr/sbin<% if node['languages']['ruby']['gem_bin'] %>:<%= node['languages']['ruby']['gem_bin'] %>/bin<% end -%>
KIBANA="<%= node['logstash']['basedir'] %>/kibana/current/kibana-daemon.rb"

test -x "$KIBANA" || exit 0

exec su -c "$KIBANA $1" "<%= node['logstash']['user'] %>"
